# script/backup.service
[Unit]
Description=Borg Backup

[Service]
Type=simple
Nice=19
ExecStart=/home/vagrant/backup.sh
